# data/packages/glacies-etc/etc/s6/rc.local
#!/bin/execlineb -P
# script run in parallel with default services
