# workdir/opt/base/etc/rc.local
#!/bin/rc
#

do_disable_overcommit=true
do_protected_symlinks=false
do_protected_hardlinks=true

if ($do_disable_overcommit)
	echo 2 > /proc/sys/vm/overcommit_memory

if ($do_protected_symlinks)
	echo 1 > /proc/sys/fs/protected_symlinks

if ($do_protected_hardlinks)
	echo 1 > /proc/sys/fs/protected_hardlinks
