# etc/rc.local
#!/bin/sh
#

do_disable_overcommit=1
do_protected_symlinks=0
do_protected_hardlinks=1

if $do_disable_overcommit ; then
	printf "%s\n" 2 > /proc/sys/vm/overcommit_memory
fi

if $do_protected_symlinks ; then
	printf "%s\n" 1 > /proc/sys/fs/protected_symlinks
fi

if $do_protected_hardlinks ; then
	printf "%s\n" 1 > /proc/sys/fs/protected_hardlinks
fi
