# data/packages/glacies-etc/etc/s6/sv/mount/up
#!/bin/execlineb -P
if { mount -o remount,rw / }
if { mkdir -p /dev/pts /dev/shm }
if { mount -a }
if { swapon -a }
ln -s /proc/self/fd /dev/fd
