# data/packages/glacies-etc/etc/s6/sv/modules/up
#!/bin/execlineb -P
backtick -D "" modules {
	redirfd -w 2 /dev/null
	venus-conf -t /etc/s6/rc.conf modules
}
if -t {
	importas -iu modules modules
	test -n "${modules}"
}
importas -isu modules modules
modprobe -aq $modules
