# data/packages/glacies-etc/etc/s6/sv/modules/up
#!/bin/execlineb -P
if -t { test -s /etc/modules }
pipeline { cat /etc/modules } xargs modprobe -aq
