# data/packages/glacies-etc/etc/s6/sv/tty1/run
#!/bin/execlineb -P
foreground { if { test -e /etc/issue } cat /etc/issue }
getty /dev/tty1
